# site_local/etc/crontab
#
* 0,12 * * * /home/webservice/projects/MyApps/script/batch.pl getThumbs
